# data/udev/90-HOTPLUG_display-visor.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="pkill -RTMIN+5 display-visor"
